# pkgs/flox-base-profile/profile.d/0500_python3.sh
# ============================================================================ #
#
# Setup Python3
#
# ---------------------------------------------------------------------------- #

# Only run if `python3' is in `PATH'
if command -v python3 >/dev/null; then
  PYTHONPATH="${PYTHONPATH:+$PYTHONPATH:}$(
    python3 -c 'import site; print( site.USER_SITE );'
  )";
  export PYTHONPATH;
fi


# ---------------------------------------------------------------------------- #
#
#
#
# ============================================================================ #
